# roles/journald/templates/journald-logger.cron.j2
# Restart journald logger to vaccum log
# {{ ansible_managed }}
0 3 * * * root /usr/bin/systemctl restart journald-logger.service
